# img_network/systemd/training_daemon@.service
[Unit]
Description=OCR training daemon

[Service]
Type=simple
WorkingDirectory=${working_dir}
ExecStart=${bin_path} %i ${training_set}